# roles/rpi2/templates/hdd.mount.j2
[Unit]
Description=hdd

[Mount]
What={{ hdd_src }}
Where={{ hdd_dest }}
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target
