# roles/rpi3/files/media-HDD.mount
[Unit]
Description=hdd

[Mount]
What=/dev/sda1
Where=/media/HDD
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target
